# build_distribution_package/linux/robbojunior-x32/etc/xdg/autostart/junior-shortcut-added.desktop
[Desktop Entry]
Name=Robbo Junior
Icon=junior
Exec=/opt/junior/create_shortcut.sh
Terminal=false
Type=Application
StartupNotify=false
Name[ru]=junior-shortcut-added.desktop
